# modules/OpenLDAPBis/ldap.cron
45       22      *       *       *	root       backupLDAP.sh >> /var/log/backupLDAP-ok.log 2> /var/log/backupLDAP-erro.log 
